# crontab
0 5 * * * root python3 /etc/repcal/genRSS.python
#blank line 